# fireberry.service
[Unit]
Description=fireberry
After=multi-user.target

[Service]
Type=idle
ExecStart=/usr/local/bin/fireberry

[Install]
WantedBy=multi-user.target
